# deployment/systemd/11of-frontend.service
[Unit]
Description=11of Frontend Service
After=network.target 11of-backend.service
Wants=11of-backend.service

[Service]
Type=simple
User=11of-app
Group=11of-app
WorkingDirectory=/opt/11of/frontend/current
Environment="NODE_ENV=production"
Environment="PORT=8080"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/11of/frontend/current/.env.production
ExecStart=/usr/bin/node /opt/11of/frontend/current/node_modules/.bin/next start -p 8080
Restart=on-failure
RestartSec=10
StandardOutput=append:/opt/11of/frontend/logs/app.log
StandardError=append:/opt/11of/frontend/logs/error.log

# Relaxed security for Node.js
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/11of/frontend

[Install]
WantedBy=multi-user.target
